# scripts/com.hermes.omlx-controlled.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hermes.omlx-controlled</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/xqdwww/Workspace/AI_Core/hermes-agent/scripts/omlx_start_controlled.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/xqdwww/.omlx/logs/hermes-omlx-launchagent.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/xqdwww/.omlx/logs/hermes-omlx-launchagent.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OMLX_MEMORY_GUARD_GB</key>
    <string>120</string>
    <key>OMLX_MODEL_DIR</key>
    <string>/Users/xqdwww/Workspace/AI_Core/mlx_models</string>
    <key>OMLX_BASE_PATH</key>
    <string>/Users/xqdwww/.omlx</string>
    <key>OMLX_LOG_DIR</key>
    <string>/Users/xqdwww/.omlx/logs</string>
  </dict>
</dict>
</plist>
